# 2_raid/raid-part1.mount
[Unit]
Description=Part 1

[Mount]
What=/dev/md127p1
Where=/raid/part1
Type=ext4
Options=defaults

[Install]
WantedBy=multi-user.target