# 2_raid/raid-part5.mount
[Unit]
Description=Part 5

[Mount]
What=/dev/md127p5
Where=/raid/part5
Type=ext4
Options=defaults

[Install]
WantedBy=multi-user.target